# nginx_rp/s6/certbot/finish
#!/bin/sh
echo "certbot finish script"
# RUN_EXIT_CODE=$1
# if [ $RUN_EXIT_CODE -ne 0 ]; then
#     echo $RUN_EXIT_CODE
# fi
exit 125
