# nginx_rp/s6/certbot/finish
#!/bin/sh
s6-svc -O /s6/certbot
echo "certbot first run complete"
