# nginx_rp/s6/certbot/finish
#!/bin/sh
echo "certbot first run complete"
echo $?
